# consul/client/consul.conf
[program:consul]
command=/usr/local/sbin/consul agent -join=172.17.42.1 -data-dir=/tmp/consul
autostart=true
autorestart=true