# consul/server/consul.conf
[program:consul]
command=/usr/local/sbin/consul agent -server -join=172.17.42.1 -data-dir=/tmp/consul -config-dir=/consul/config
autostart=true
autorestart=true